# etc/uci-defaults/luci-theme-mytheme
#!/bin/sh
	uci batch <<-EOF
	        set luci.themes.MyTheme=/luci-static/mytheme
	        commit luci
	EOF
exit 0